# desk/deploy/desk-cycle.timer
[Unit]
Description=Desk reflective cycle — autopsy, briefs, self-mod gate, digest (daily)

[Timer]
# Once a day after the morning settlement, when yesterday's highs are resolved.
OnCalendar=*-*-* 14:15:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
