# desk/deploy/desk-cycle.service
[Unit]
Description=Desk reflective cycle (fake money only — learning loop)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/polymarket
# The LLM backend is chosen by env. Leave DESK_LLM=mock for a zero-cost dry loop;
# set DESK_LLM=claude + ANTHROPIC_API_KEY (in the EnvironmentFile) to go live.
EnvironmentFile=-%h/polymarket/desk/.env
ExecStart=/usr/bin/python3 -m desk.run_cycle
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=%h/polymarket
