# desk/deploy/papertrader-scan.service
[Unit]
Description=PaperTrader reactive scan (fake money only — scan + settle)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/polymarket/papertrader
ExecStart=/usr/bin/python3 papertrader.py scan
ExecStartPost=/usr/bin/python3 papertrader.py settle
ExecStartPost=/usr/bin/python3 -m desk.export_state
ExecStartPost=/bin/bash -c 'cd %h/polymarket && git add papertrader/data desk/dashboard_state.json desk/memory/lessons && git commit -m "auto(scan): $(date -u +%%FT%%TZ)" || true && git push'
# hardening (defense in depth on the VPS)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=%h/polymarket
ProtectHome=read-only
