# desk/deploy/papertrader-scan.timer
[Unit]
Description=PaperTrader reactive scan — every 30 minutes

[Timer]
OnCalendar=*-*-* *:00,30:00 UTC
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
